# 4-nginx-reverse-proxy-fat-jar/chef-repo/cookbooks/devopsdays/templates/default/fat_jar.service.erb
[Unit]
Description=<%= @serviceDescription %>
After=syslog.target
After=network.target

[Service]
Type=simple
User=<%= @user %>
ExecStart=/apps/java/bin/java <%= @javaOpts%> -Dfile.encoding=UTF-8 -jar <%= @jarPath %> <%= @progArgs%>
ExecStop=/bin/kill -TERM $MAINPID
TimeoutSec=300

[Install]
WantedBy=multi-user.target